# services/system/health-server.service
[Unit]
Description=Enhanced Health Check API Server
After=network.target openclaw-gateway.service
Wants=openclaw-gateway.service

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/health_server.py
Restart=always
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
